# openrc/reverssh
#!/sbin/openrc-run

supervisor=supervise-daemon

description="Reversing SSH tarpit"
command="/usr/local/bin/reverssh"
command_args="-c /etc/reverssh/reverssh.toml"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/${RC_SVCNAME}.log"
error_log="/var/log/${RC_SVCNAME}.log"

depend() {
	need net
	after firewall
}